# lib/systemd/system/stream.service
[Service]
ExecStart=stream
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
